# src/program-redis.conf
[program:redis]
command = /usr/local/bin/launch-redis.sh
autorestart = true
